# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf
